# deployment/home-dashboard-backend.service
# ─────────────────────────────────────────────────────────────────────────────
# systemd unit for the Home Dashboard Rust backend
# ─────────────────────────────────────────────────────────────────────────────
# Install:
#   sudo cp home-dashboard-backend.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now home-dashboard-backend
#
# Check status:
#   systemctl status home-dashboard-backend
#   journalctl -u home-dashboard-backend -f
# ─────────────────────────────────────────────────────────────────────────────

[Unit]
Description=Home Dashboard Backend — System Metrics & Supabase Keep-Alive
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=adhaghani
Group=adhaghani
ExecStart=/usr/local/bin/home-dashboard-backend

# Always restart on crash (critical for keep-alive)
Restart=always
RestartSec=5

# Prevent thrashing: stop restarting if it fails 3 times in 60 seconds
StartLimitInterval=60
StartLimitBurst=3

# ── Environment variables ──────────────────────────────────────────────────
Environment=SUPABASE_URL=https://tytyvuiftcfjzosglsbm.supabase.co
Environment=SUPABASE_ANON_KEY=sb_publishable_Ya6Hz0s1O_nY92JrlEV4iw_fgr_kj6Q
Environment=RUST_LOG=home_dashboard_backend=info
StandardOutput=journal
StandardError=journal
SyslogIdentifier=home-dashboard-backend

# Security hardening (relaxed — need tailscale socket and pi.hole host resolution)
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/tmp
ReadOnlyPaths=/sys/class/thermal /home/adhaghani/.nanobot/config.json

[Install]
WantedBy=multi-user.target
